# recorder/ai.reticle.meeting-recorder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>ai.openclaw.meeting-recorder</string>

	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/meeting-recorder</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>

	<key>ThrottleInterval</key>
	<integer>30</integer>

	<key>StandardOutPath</key>
	<string>/Users/alexandervyhmeister/GitHub/claudia/meeting-recorder.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/alexandervyhmeister/GitHub/claudia/meeting-recorder-error.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/alexandervyhmeister</string>
	</dict>
</dict>
</plist>
